# systemd-networkd-wait-online.service.d/01-wait-online.conf
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online -i eth0 -i eth1  -i virbr0  -i virbr1 -i wlan0 -i wlan1 --any
TimeoutSec=10
